# etc/.config/autostart/notify-server.desktop
[Desktop Entry]
Exec=/home/tm/.bin/scripts/notify-server.py
Name=notify-server
Type=Application
